# files/podget.timer
[Unit]
Description=Run podget daily

[Timer]
OnCalendar=*-*-* 10:05:00
Persistent=true

[Install]
WantedBy=timers.target
