# files/podget.service
[Unit]
Description=Run podget

[Service]
Type=oneshot
ExecStart=/usr/bin/podget -s
